# .config/autostart/picom.desktop
[Desktop Entry]
Type=Application
Name=picom
Exec=picom -b
StartupNotify=false
Terminal=false
