# .config/autostart/udiskie.desktop
[Desktop Entry]
Name=udiskie
Type=Application
Exec=udiskie --automount
NoDisplay=true
