# .config/autostart/udiskie.desktop
[Desktop Entry]
Type=Application
Name=udiskie
Comment=Automounter for removable media
Exec=udiskie --automount
Terminal=false
NoDisplay=true
